# deploy/todoist-mcp.service
[Unit]
Description=Todoist MCP Server
After=docker.service
Requires=docker.service

[Service]
Type=simple
Restart=always
RestartSec=10
User=docker
Group=docker

# Pull the latest image before starting
ExecStartPre=/usr/bin/docker pull todoist-mcp:latest

# Stop and remove any existing container
ExecStartPre=-/usr/bin/docker stop todoist-mcp
ExecStartPre=-/usr/bin/docker rm todoist-mcp

# Run the container with remote server
ExecStart=/usr/bin/docker run \
    --name todoist-mcp \
    --rm \
    -p 8765:8765 \
    -v /opt/todoist-mcp/data:/root/.todoist-mcp \
    -v /opt/todoist-mcp/.env:/app/.env:ro \
    --env-file /opt/todoist-mcp/.env \
    todoist-mcp:latest \
    python -m src.remote_server

# Stop the container gracefully
ExecStop=/usr/bin/docker stop todoist-mcp

[Install]
WantedBy=multi-user.target